# etc/init/nightcap.conf
description "nightcap"
author      "Daniel Rice"

start on local-filesystems
stop on shutdown

respawn

exec /mnt/data/projects/nightcap/nightcap
